# etcfs/profile.d/eduldap.sh
#!/usr/bin/sh

touch /tmp/wtf
#export PATH="/etc/eduldap/bin:${PATH}"
